# deploy/civilpdf-backend.service
[Unit]
Description=CivilPDF-DX Backend (FastAPI / uvicorn)
Documentation=https://github.com/kensan/CivilPDF-DX
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
WorkingDirectory=/home/kensan/Projects/CivilPDF-DX/src/console/backend

# Load secrets from the env file (not checked into git)
EnvironmentFile=/home/kensan/Projects/CivilPDF-DX/deploy/civilpdf.env

# Run DB migration then start the server
ExecStartPre=/bin/bash -c '/home/kensan/.local/bin/alembic upgrade head'
ExecStart=/home/kensan/.local/bin/uvicorn main:app \
    --host 0.0.0.0 \
    --port 8000 \
    --workers 2 \
    --log-level info

Restart=on-failure
RestartSec=5s

# Upload directory must exist
RuntimeDirectory=civildx
RuntimeDirectoryMode=0755

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
